# deploy/ai-issue-runner/clankerdev-ai-issue-runner.service
[Unit]
Description=Clankerdev Codex issue-to-PR runner
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=HOME=/root
Environment=XDG_CONFIG_HOME=/root/.config
Environment=CODEX_HOME=/root/.codex
Environment="GIT_SSH_COMMAND=ssh -o BatchMode=yes"
EnvironmentFile=-/etc/clankerdev-ai/runner.env
ExecStart=/usr/local/bin/clankerdev-ai-issue-runner
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
